# roles/mailserver/files/httpd/md-message-installed2@.service
[Service]
Type=exec
ExecStart=/usr/local/libexec/httpd/md-message __installed2 %I
Environment=MD_MESSAGE_LOGGER=journal
