# roles/ipa-server/files/ipa-pki-crl-fetch.service
[Unit]
Description=IPA PKI CRL fetch

[Service]
Type=exec
ExecStart=/usr/local/libexec/ipa-pki-crl-fetch
