# roles/ipa-server/files/pki-log-cleanup.timer
[Unit]
Description=pki log file cleanup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=14400

[Install]
WantedBy=timers.target
